# opt/unipiap/systemd/unipiap.target
[Unit]
Description=Wifi AP for UniPi configuration
